# update-pkgs.service
[Unit]
Description=Hourly Package Update Script

[Service]
Type=oneshot
ExecStart=/home/amefan/updates.sh

[Install]
WantedBy=multi-user.target
