# update-pkgs.service
[Unit]
Description=Hourly Package Update Script

[Service]
Type=oneshot
ExecStart=$HOME/updates.sh

[Install]
WantedBy=multi-user.target
